# showtime-server.service
[Unit]
Description=Showtime Web Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/showtime-server
Restart=always
RestartSec=3
KillMode=mixed

[Install]
WantedBy=multi-user.target
